# demo/payload/supervisor/mysql.conf
[program:mysql]
autorestart=true
command=/usr/bin/mysqld_safe
