# demo/payload/supervisor/redis.conf
[program:redis]
autorestart=true
command=/usr/bin/redis-server
