# supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server
autorestart=true
user=root
stdout_logfile=/var/log/redis/stdout.log
stderr_logfile=/var/log/redis/stderr.log